# uploadmails/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*     *       *       *       *       fetch flock -x -n /tmp/fetch.lock /srv/fetch/grab.sh /var/fetch

